# llama-server.service
[Unit]
Description=llama.cpp Server Daemon
After=network.target

[Service]
Type=simple
User=%U
Group=%g
WorkingDirectory=%h/git/llama.cpp
EnvironmentFile=-%h/git/llama.cpp-daemon/.env
ExecStart=%h/git/llama.cpp/build/bin/llama-server \
  -m ${MODEL_PATH} \
  --host ${HOST} --port ${PORT} \
  -ngl ${NGL_LEVEL} --ctx-size ${CONTEXT_SIZE} ${JINJA_ENABLED}
Restart=${RESTART_MODE}
RestartSec=${RESTART_SECONDS}
StandardOutput=journal
StandardError=journal
SyslogIdentifier=${SERVICE_NAME}

[Install]
WantedBy=multi-user.target
